# security-system-web.service
[Unit]
Description=Riftech Cam Security Web Interface
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/riftech-cam-security
Environment="PATH=/home/YOUR_USERNAME/riftech-cam-security/venv/bin"
ExecStart=/home/YOUR_USERNAME/riftech-cam-security/start_web.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
